# t1-payout-app-root.service
[Unit]
Description=T1 Solana Payout App Root Static Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/t1-payout-app-root
ExecStart=/usr/bin/python3 -m http.server 8091 --directory %h/t1-payout-app-root
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
